# deploy/the-station-agent.service
[Unit]
Description=The Station Agent - Flask Backend Service
After=network.target mongodb.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/root/thestationagentblue.fun

# ⭐ 加载环境变量（关键配置）
EnvironmentFile=/root/thestationagentblue.fun/.env

# 虚拟环境 Python 解释器路径
Environment="PATH=/root/thestationagentblue.fun/.venv/bin"

# Gunicorn 启动命令
# -w 3: 3 个工作进程
# -b 127.0.0.1:8000: 绑定本地 8000 端口（由 Nginx 反向代理）
# --timeout 120: 请求超时时间 120 秒
# --access-logfile -: 访问日志输出到 systemd journal
# --error-logfile -: 错误日志输出到 systemd journal
ExecStart=/root/thestationagentblue.fun/.venv/bin/gunicorn \
    -w 3 \
    -b 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    app:app

# 自动重启策略
Restart=always
RestartSec=10

# 进程管理
StandardOutput=journal
StandardError=journal
SyslogIdentifier=the-station-agent

[Install]
WantedBy=multi-user.target
